# docker_files/supervisord.conf
[supervisord]
nodaemon=true

[program:trytond]
command=/usr/local/bin/trytond -c /etc/trytond.conf

[program:nereid]
command=python application_example.py
directory=/opt/nereid-webshop/web
